# cron/conf/www-data.crontab
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
5       *       *       *       *       /usr/bin/env php /var/www/html/console core:archive --url=https://${MATOMO_FQDN}/
